# irc/services.d/sshd/run
#!/usr/bin/with-contenv sh

# Launch ssh server, stay in foreground
/usr/sbin/sshd -D $sshd_args

