# forked-daapd/services.d/dbus/run
#!/bin/sh

mkdir -p /var/run/dbus
chown -R messagebus:messagebus /var/run/dbus
dbus-uuidgen --ensure
dbus-daemon --nofork --system --nopidfile
