# forked-daapd/services.d/avahi/run
#!/usr/bin/execlineb -P
avahi-daemon --no-chroot
